# gitops_sync/rootfs/etc/services.d/gitops-sync/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting Home Assistant GitOps Sync"
exec python3 /opt/gitops-sync/main.py
